# com.openclaw.daily-content.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.daily-content</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/openclaw/trading-framework/pipelines/run_daily_content.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/openclaw/trading-framework</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/openclaw/trading-framework/pipelines/data/content/daily_content.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/openclaw/trading-framework/pipelines/data/content/daily_content.log</string>
</dict>
</plist>
