# com.openclaw.ceowatcher-reader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.ceowatcher-reader</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/openclaw/trading-framework/pipelines/ceowatcher_reader.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/openclaw/trading-framework</string>

    <!-- Run daily at 8:30 AM ET (before market open) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/openclaw/trading-framework/logs/ceowatcher_reader.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/openclaw/trading-framework/logs/ceowatcher_reader.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/openclaw/trading-framework</string>
    </dict>
</dict>
</plist>
