# com.openclaw.trial-emails.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.trial-emails</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/openclaw/trading-framework/pipelines/trial_emails.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/openclaw/trading-framework</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>14</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/openclaw/trading-framework/logs/trial_emails.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/openclaw/trading-framework/logs/trial_emails.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/openclaw/trading-framework</string>
    </dict>
</dict>
</plist>
