# com.openclaw.breaking-signal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.breaking-signal</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/openclaw/trading-framework/pipelines/run_breaking_signal.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/openclaw/trading-framework</string>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>StandardOutPath</key>
    <string>/Users/openclaw/trading-framework/pipelines/data/content/breaking_signal.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/openclaw/trading-framework/pipelines/data/content/breaking_signal.log</string>
</dict>
</plist>
