# com.openclaw.backfill-returns.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.backfill-returns</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/openclaw/trading-framework/pipelines/insider_study/backfill_returns.py</string>
        <string>--max-download</string>
        <string>500</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/openclaw/trading-framework</string>

    <!-- Run daily at 5 AM (after markets close, before morning review) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/openclaw/trading-framework/logs/backfill_returns.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/openclaw/trading-framework/logs/backfill_returns.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/openclaw/trading-framework</string>
    </dict>
</dict>
</plist>
