# scripts/launchd/com.openclaw.post-deploy-audit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.post-deploy-audit</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/derekg/trading-framework/scripts/post_deploy_audit.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/derekg/trading-framework</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/derekg/trading-framework</string>
    </dict>
    <!-- Originally one-shot Monday 2026-05-11 14:00 PT (post-deploy audit).
         After the Monday run revealed daily-prices being orphaned from the
         morning refresh chain (now fixed), extending to fire each weekday
         14:00 PT this week (Tue 5/12, Wed 5/13, Thu 5/14, Fri 5/15) for
         daily confirmation that the fix holds. Each StartCalendarInterval
         entry is a separate date — launchd OR's across them. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Year</key><integer>2026</integer>
            <key>Month</key><integer>5</integer>
            <key>Day</key><integer>12</integer>
            <key>Hour</key><integer>14</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Year</key><integer>2026</integer>
            <key>Month</key><integer>5</integer>
            <key>Day</key><integer>13</integer>
            <key>Hour</key><integer>14</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Year</key><integer>2026</integer>
            <key>Month</key><integer>5</integer>
            <key>Day</key><integer>14</integer>
            <key>Hour</key><integer>14</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Year</key><integer>2026</integer>
            <key>Month</key><integer>5</integer>
            <key>Day</key><integer>15</integer>
            <key>Hour</key><integer>14</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/derekg/trading-framework/logs/post-deploy-audit-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/derekg/trading-framework/logs/post-deploy-audit-launchd.log</string>
</dict>
</plist>
