# scripts/launchd/com.openclaw.drift-detector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.drift-detector</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/derekg/trading-framework/scripts/drift_detector.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/derekg/trading-framework</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/derekg/trading-framework</string>
    </dict>
    <!-- Run daily at 08:00 PT — after strategy-simulator (07:00) and after
         cw_runner morning cycle. Reads sim_portfolio (the simulator's
         current view) and paper_trades / strategy_portfolio (the runner's
         actual fills) and records every divergence in strategy_drift_audit. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>8</integer>
        <key>Minute</key><integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/derekg/trading-framework/logs/drift-detector.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/derekg/trading-framework/logs/drift-detector.log</string>
</dict>
</plist>
