# scripts/launchd/com.openclaw.freshness-probe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.freshness-probe</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/derekg/trading-framework/scripts/freshness_probe.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/derekg/trading-framework</string>
    <!-- Every 30 minutes. Cheap (single PG read per contract); the alert
         layer fires only on transitions, so frequent runs don't spam. -->
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/derekg/trading-framework/logs/freshness-probe.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/derekg/trading-framework/logs/freshness-probe.log</string>
    <key>ProcessType</key>
    <string>Background</string>
    <key>LowPriorityIO</key>
    <true/>
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
