# strategies/insider_catalog/com.openclaw.refresh-features.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.refresh-features</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/derekg/trading-framework/strategies/insider_catalog/refresh_features_daily.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/derekg/trading-framework</string>
    <!-- Runs at 06:00 PT (09:00 ET) every weekday. Studio is in PT — launchd
         StartCalendarInterval uses local time. The cw_runner strategies
         wake at 09:25 ET (= 06:25 PT) — refresh-features must complete BEFORE
         that wake or the runners halt at preflight with R-002 (stale data).
         The chain takes ~30 sec for the 30-day --since window, so 25 min
         margin is comfortable.

         Previous schedule (09:30 PT = 12:30 ET) ran 3 HOURS AFTER the runner
         wake; strategies halted every market morning until refresh finished
         at 12:30+ ET. Phase 2 P1 fix: align so refresh completes BEFORE wake. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>6</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>6</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>6</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>6</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>6</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/derekg/trading-framework/logs/refresh-features.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/derekg/trading-framework/logs/refresh-features.log</string>
    <key>ProcessType</key>
    <string>Background</string>
    <key>LowPriorityIO</key>
    <true/>
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
