# dataplane/deploy/com.openclaw.dagster-webserver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.dagster-webserver</string>

  <!-- Binds to the Tailscale IP only: the Dagster UI has no auth, so it
       must be reachable from Derek's tailnet devices but not the LAN.
       If tailscaled isn't up yet at load, the bind fails and KeepAlive
       retries every 30s until it succeeds — self-healing on boot. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>set -a; . /Users/derekg/trading-framework/.env 2>/dev/null || true; set +a; exec /Users/derekg/dataplane_venv/bin/dagster-webserver -h 100.78.9.66 -p 3000</string>
  </array>

  <!-- cwd must be DAGSTER_HOME: webserver reads workspace.yaml from cwd. -->
  <key>WorkingDirectory</key>
  <string>/Users/derekg/dataplane_dagster_home</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>DAGSTER_HOME</key>
    <string>/Users/derekg/dataplane_dagster_home</string>
    <key>DBT_PROFILES_DIR</key>
    <string>/Users/derekg/trading-framework/dataplane/dbt_project</string>
    <key>PATH</key>
    <string>/Users/derekg/dataplane_venv/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/Users/derekg/dataplane_dagster_home/logs/webserver.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/derekg/dataplane_dagster_home/logs/webserver.err</string>
</dict>
</plist>
